# deploy/systemd-backend.service.example
# Systemd сервис для Django Backend
# Скопируйте в /etc/systemd/system/sofa-backend.service
# Замените /home/deploy на путь к вашему проекту и deploy на вашего пользователя

[Unit]
Description=Sofa Marketplace Django Backend
After=network.target postgresql.service

[Service]
Type=notify
User=deploy
Group=deploy
WorkingDirectory=/home/deploy/sofa-marketplace/backend
Environment="PATH=/home/deploy/sofa-marketplace/backend/venv/bin"
ExecStart=/home/deploy/sofa-marketplace/backend/venv/bin/gunicorn \
    --config gunicorn_config.py \
    config.wsgi:application
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target

